# configs/supervisor.ini
[program:co2info]
command = /home/pi/co2/server.py --logging=info
user = root
redirect_stderr = true
stdout_logfile = /home/pi/co2/logs/server.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 1
directory=/home/pi/co2/

